# dist/DEBIAN/postinst
#!/bin/sh
set -e
if [ "$1" = "configure" ] || [ "$1" = "abort-upgrade" ] || [ "$1" = "abort-deconfigure" ] || [ "$1" = "abort-remove" ] ; then
	if [ -d /run/systemd/system ]; then
		systemctl --system daemon-reload >/dev/null || true
        	deb-systemd-helper unmask 'p2psec.service' >/dev/null || true
		deb-systemd-helper enable 'p2psec.service' >/dev/null || true
	fi
fi
